# meteor.sv.conf
[program:meteor]

; Set the following environment variables for your Meteor application
;
; 1. MONGO_URL (required, the address to your mongo deployment) e.g. "mongodb://user:password@your-mongo-host.com:27017/database_name"
; 2. ROOT_URL (required, the public access URL for your app) e.g. "http://www.sample.com"
; 3. PORT e.g. "80" (required, must match exposed port in Dockerfile, 8000 by default)
; 4. USE_JSESSIONID (required for load-balanced IE support) set to "1"
; 5. MONGO_OPLOG_URL (optional, use for oplog-based horizontal scaling) e.g. "mongodb://user:password@your-mongo-host.com:27017/local?authSource=database_name"
;
; In the form of:
; environment=KEY1="VALUE1",KEY2="VALUE2"

environment=MONGO_URL="",PORT="80",ROOT_URL="",USE_JSESSIONID="1"

; DO NOT EDIT BELOW

command=node /meteor/main.js
redirect_stderr=true
startsecs=5
